# deploy/gcp/certbot-ip-renew.service
[Unit]
Description=Renew short-lived Let's Encrypt IP certificates
After=network-online.target nginx.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/certbot-ip/bin/certbot renew --quiet
ExecStartPost=/usr/sbin/nginx -t
ExecStartPost=/bin/systemctl reload nginx
